# root_settings/etc/udev/rules.d/95-monitor-hotplug.rules
ACTION=="change", SUBSYSTEM=="drm", RUN+="/media/files/workspace/dotrc/other_files/xrandr.sh udev"
