# home_settings/.config/systemd/user/default_user_gui.target
[Unit]
Requires=i3wm.service
# Requires=thunderbird.service
