# home_settings/.config/systemd/user/generate_configs.service
[Service]
Type=oneshot
ExecStart=/bin/bash -c ~/os_settings/other_files/generate_configs.sh
