# other_files/old_root_settings/etc/systemd/system/default_system_gui.target
[Unit]
Requires=default_system_gui.service
