# root_settings/etc/systemd/system/default_system_gui.target
[Unit]
Requires=multi-user.target
Requires=default_system_gui.service
Conflicts=rescue.target
Conflicts=graphical.target
