# root_settings/etc/systemd/system/noob.target
# A target to run startx from user 'noob'

[Unit]
Description=startx noob
Requires=multi-user.target
Conflicts=rescue.target

